# templates/etc/profile.d/umask.sh.erb
### managed by secc_os_linux ###
umask <%= profile_umask %>
