# files/etc/profile.d/umask.sh
### managed by secc_os_linux ###
umask 027

